# qubes-integration/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (WebSocket + Dashboard)
After=network-online.target openclaw-cursor-proxy.service
Wants=network-online.target
Requires=openclaw-cursor-proxy.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStartPre=/bin/sh -c "fuser -k 18789/tcp 2>/dev/null || true"
ExecStart=%h/.npm-global/bin/openclaw gateway --bind loopback
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10
Environment=PATH=%h/.local/bin:%h/.npm-global/bin:/usr/local/bin:/usr/bin
Environment=HOME=%h
MemoryMax=1G
TasksMax=128
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-gateway

[Install]
WantedBy=default.target
